# packages/raspberrypi-hdmi/raspberrypi-hdmi-off.service
[Unit]
Description=Disable HDMI output on Raspberry Pi
After=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=/opt/vc/bin/tvservice --off
ExecStop=/opt/vc/bin/tvservice --preferred
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
