# packages/raspberrypi-hdmi/raspberrypi-hdmi-off.timer
[Unit]
Description=Disable HDMI output on Raspberry Pi after 15mins after boot

[Timer]
OnBootSec=15min

[Install]
WantedBy=timers.target
