# packages/tailscale/tailscale.install
post_install() {
	systemctl daemon-reload

	echo -e "\n\e[1m Tailscale $1\e[0m\n"
	cat <<EOF
To start using Tailscale, run as root:

	systemctl enable --now tailscaled
	tailscale up

EOF
}

post_upgrade() {
	systemctl daemon-reload
}

pre_remove() {
	systemctl disable --now tailscaled
}
