# packages/tailscale-pikvm/tailscale-pikvm.install
pre_upgrade() {
	service_enabled=$(systemctl is-enabled tailscaled || true)
}

post_install() {
	post_upgrade
}

post_upgrade() {
	mkdir -p /var/lib/tailscale || true
	chmod 700 /var/lib/tailscale || true
	if [ "$service_enabled" == "enabled" ]; then
		systemctl enable tailscaled || true
	fi
}
