# packages/kvmd-oled/kvmd-oled.service
[Unit]
Description=Pi-KVM - Small OLED daemon
After=systemd-modules-load.service
ConditionPathExists=/dev/i2c-1

[Service]
Type=simple
Restart=always
RestartSec=3
ExecStartPre=/usr/bin/kvmd-oled --height=32 --interval=3 --clear-on-exit --image=/usr/share/kvmd-oled/hello.ppm
ExecStart=/usr/bin/kvmd-oled --height=32 --offset-y=6 --clear-on-exit
TimeoutStopSec=3

[Install]
WantedBy=multi-user.target
